# src/main/dist/bin/ewdb.service
[Unit]
Description=EasyWallet HyperSQL DB

[Service]
ExecStart=/opt/EasyWallet/bin/ewdb
User=easywallet

[Install]
WantedBy=multi-user.target
